# root/etc/s6-overlay/s6-rc.d/wings/run
#!/command/with-contenv bash

wings --config $DATADIR/config.yml