# root/etc/services.d/wings/run
#!/usr/bin/with-contenv bash

# Run Wings, with specified data directory
wings --config $DATADIR/config.yml